# test-apps/libLogHelper/cmake_install_scripts/rc.local.in
#!/bin/bash

sudo @HOME_DIR@/atBoot.sh &

exit 0